# www/etc/systemd/tftp-server.socket
[Unit]
Description=Tftp Server
requires=tftp.socket


[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target
